# templates/victoria-vmauth.service.erb
[Unit]
Description=Simple auth proxy, router and load balancer for VictoriaMetrics
Documentation=https://victoriametrics.github.io/
After=network.target

[Service]
Type=simple
User=<%= @user %>
LimitNOFILE=65536
LimitNPROC=32000
StartLimitBurst=5
StartLimitInterval=0
Restart=on-failure
RestartSec=5
WorkingDirectory=<%= @storage_dir %>
EnvironmentFile=/etc/default/<%= @service_name %>
ExecStart=/usr/bin/vmauth $ARGS
ExecReload=/bin/kill -SIGHUP $MAINPID

[Install]
WantedBy=multi-user.target
