# templates/victoria-vmalert.service.erb
[Unit]
Description=Executes a list of the given alerting or recording rules against configured -datasource.url compatible with Prometheus HTTP API
Documentation=https://victoriametrics.github.io/
After=network.target

[Service]
Type=simple
User=<%= @user %>
LimitNOFILE=65536
LimitNPROC=32000
StartLimitBurst=5
StartLimitInterval=0
Restart=on-failure
RestartSec=5
WorkingDirectory=<%= @storage_dir %>
EnvironmentFile=/etc/default/<%= @service_name %>
ExecStart=/usr/bin/vmalert $ARGS
ExecReload=/bin/kill -SIGHUP $MAINPID

[Install]
WantedBy=multi-user.target
